# templates/aries/dlink.2013gui/dir883l/udev.rules/50-wan-link.rules
SUBSYSTEM=="linkreport", ACTION=="LINKUP", RUN+="/usr/sbin/usockc /var/gpio_ctrl LINKUP"
SUBSYSTEM=="linkreport", ACTION=="LINKDOWN", RUN+="/usr/sbin/usockc /var/gpio_ctrl LINKDOWN"
