# templates/aries/dlink.2013gui/dir883l/udev.rules/80-firmware.rules
# firmware-class requests, copies files into the kernel
SUBSYSTEM=="firmware", ACTION=="add", RUN+="/etc/scripts/firmware.sh"
